# ops/systemd/moneyflow-backup.timer
[Unit]
Description=Run the MoneyFlow encrypted backup daily

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=900
AccuracySec=1min
Unit=moneyflow-backup.service

[Install]
WantedBy=timers.target
