# ops/systemd/moneyflow-backup.service
[Unit]
Description=MoneyFlow encrypted PostgreSQL backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
EnvironmentFile=/opt/moneyflow/.env
ExecStart=/opt/moneyflow/ops/backup.sh
UMask=0077
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX
SystemCallArchitectures=native
ReadWritePaths=/var/lib/moneyflow-backups
